# crontab
@daily ~wilbowma/procmail/daily-sa-learn.sh